# etc/crontab
13 3 * * * root /usr/local/sbin/daily-shutdown.sh
